# budge/systemd/budge-categorize.timer.in
[Unit]
Description=budge: daily categorize timer (after fetch)

[Timer]
OnCalendar=@SCHEDULE_CATEGORIZE@
Persistent=true

[Install]
WantedBy=timers.target
