# budge/systemd/budge-review-nudge.service.in
# budge-review-nudge — weekly "review ready — N transactions pending" (7.7)
[Unit]
Description=budge: weekly review-ready nudge via OpenClaw
After=budge-categorize.service

[Service]
Type=oneshot
User=@USER@
Environment=BUDGE_CONFIG_DIR=@CONFIG_DIR@
Environment=LC_ALL=C.UTF-8
ExecStart=@BUDGE@ notify --review-nudge
SyslogIdentifier=budge-review-nudge
