# budge/systemd/budge-push.timer.in
[Unit]
Description=budge: daily push timer (after fetch + categorize)

[Timer]
OnCalendar=@SCHEDULE_PUSH@
Persistent=true

[Install]
WantedBy=timers.target
