# budge/systemd/hledger-web.service.in
# hledger-web — the official hledger web UI, served VIEW-ONLY (PRD: budge
# never grants any UI write access to the journal; edits flow through the
# review pipeline or the operator's editor).
#
# Like Paisa, this has no authentication — keep it LAN-only.
[Unit]
Description=hledger-web (official web UI, view-only) for the budge journal
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=@USER@
Environment=LC_ALL=C.UTF-8
# --allow=view makes it read-only (flag replaced --capabilities in 1.26)
ExecStart=@HLEDGER_WEB@ --serve --host=0.0.0.0 --port=5000 \
    --allow=view -f @REPO@/main.journal
Restart=on-failure
SyslogIdentifier=hledger-web

[Install]
WantedBy=multi-user.target
