# systemd/budge-fetch.timer.in
[Unit]
Description=budge: daily fetch timer

[Timer]
OnCalendar=@SCHEDULE_FETCH@
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
